# app/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:ring-gdrive]
command=/opt/ring/run/start.sh